# raspberry/auto-wifi.service
# 自动WiFi连接服务配置文件
# 创建目的：用于在系统启动时自动连接可用的WiFi网络
# 使用方法：
# 0. 创建脚本：sudo nano /etc/systemd/system/auto-wifi.service
# 1. 启用服务：sudo systemctl enable auto-wifi.service
# 2. 启动服务：sudo systemctl start auto-wifi.service
# 3. 查看状态：sudo systemctl status auto-wifi.service
# 4. 查看日志：journalctl -u auto-wifi.service
# 5. 修改路径： ExecStart=/bin/bash /path/to/auto_connect.sh为真实路径
# 6.修改wifi： ssh上服务器之后，nmcli device wifi connect <SSID> password <密码>；然后再ssh一次 nmcli connection modify <连接名称> connection.autoconnect yes

[Unit]
Description=自动WiFi连接服务
After=network.target

[Service]
Type=oneshot
ExecStartPre=/bin/bash -c 'if nmcli -t -f NAME,autoconnect connection show | grep -q "yes$"; then exit 1; fi'
ExecStart=/bin/bash /path/to/auto_connect.sh
RemainAfterExit=yes
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=0
StartLimitBurst=0

[Install]
WantedBy=multi-user.target